# .config/systemd/user/dunst.service.d/override.conf
[Unit]
PartOf=

[Install]
WantedBy=default.target
